# crypto-analyzer/deploy/systemd/dimension-trader.service
[Unit]
Description=Crypto Analyzer Dimension Strategy Trader
After=network-online.target crypto-api.service crypto-collector.service
Wants=crypto-api.service

[Service]
Type=simple
User=crypto
Group=crypto
WorkingDirectory=/opt/crypto/crypto-analyzer
EnvironmentFile=/opt/crypto/crypto-analyzer/.env
ExecStart=/opt/crypto/crypto-analyzer/.venv/bin/python dimension_trader.py
Restart=always
RestartSec=20
StandardOutput=append:/opt/crypto/crypto-analyzer/logs/dimension.log
StandardError=append:/opt/crypto/crypto-analyzer/logs/dimension.log

[Install]
WantedBy=multi-user.target
